# deploy/kokuraasahi.ini

[program:kokuraasahi_mcl]
command=/data/kokuraasahi/mcl/java/bin/java -jar /data/kokuraasahi/mcl/mcl.jar
directory=/data/kokuraasahi/mcl
numprocs=1
autostart=true
autorestart=true
startretries=3
startsecs=10
redirect_stderr=true
stderr_logfile=/var/log/kokuraasahi/stdout.log
stdout_logfile=/var/log/kokuraasahi/stdout.log
user=pruto

